# homelab/roles/deutsch-app/templates/deutsch-app-deploy.service.j2
[Unit]
Description=Restart Deutsch App after build

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart deutsch-app
