# raspberry-pi/roles/deutsch-app/templates/deutsch-app-deploy.path.j2
[Unit]
Description=Watch for Deutsch App builds

[Path]
PathChanged={{ deutsch_app_dir }}/.next/BUILD_ID

[Install]
WantedBy=multi-user.target
