# raspberry-pi/roles/tailscale/templates/tailscale-serve.service.j2
[Unit]
Description=Tailscale Serve - HTTPS Proxy
Documentation=https://tailscale.com/kb/1242/tailscale-serve
After=network.target tailscaled.service opencode.service deutsch-app.service
Requires=tailscaled.service
Wants=opencode.service deutsch-app.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
  /usr/bin/tailscale serve reset && \
  /usr/bin/tailscale serve --bg --https=443 http://localhost:{{ opencode_port }} && \
  /usr/bin/tailscale serve --bg --https=443 --set-path /deutsch http://localhost:{{ deutsch_app_port }}/deutsch'
ExecStop=/usr/bin/tailscale serve reset

[Install]
WantedBy=multi-user.target
